# dentos-poe-agent/DEBIAN/postinst
#!/bin/bash

systemctl --no-block enable poed.service
systemctl --no-block start poed.service
